# com.echo.alert-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.echo.alert-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3.13</string>
        <string>/Users/konstantinkorobov/Knowledge/The Garden/02_Blueprint/Echo/karaoke-app/alert_bot.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/echo-alert-bot.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/echo-alert-bot.log</string>
</dict>
</plist>
